# ai-team-server.service
[Unit]
Description=AI Team Dashboard Backend
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/opt/ai-team-server
ExecStart=/usr/bin/node /opt/ai-team-server/server.js
Restart=always
RestartSec=5
Environment=NODE_ENV=production

[Install]
WantedBy=multi-user.target
